# ansible/resources/traefik.service
[Unit]
Description=Traefik service

[Service]
ExecStart=/usr/local/bin/traefik

[Install]
WantedBy=multi-user.target
